# headsetcontrol-tray/headsetcontrol-tray.service
[Unit]
Description=HeadsetControl Tray Widget
After=graphical-session.target

[Service]
ExecStart=/usr/bin/headsetcontrol-tray
Restart=on-failure
Type=simple

[Install]
WantedBy=graphical-session.target
